# .config/monitor-system/monitor-battery.service
[Unit]
Description=Monitor Battery deamon

[Service]
User=ianni
ExecStart=/bin/monitor-battery

[Install]
WantedBy=multi-user.target
